# pi-sampler.service
[Unit]
Description=Raspberry Pi MIDI Sampler
After=sound.target network.target
Wants=sound.target

[Service]
Type=simple
User=pi
Group=audio
WorkingDirectory=/home/pi/rpi-sampler
ExecStart=/home/pi/rpi-sampler/sampler
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Real-time priority settings
LimitRTPRIO=99
LimitMEMLOCK=infinity

# Environment
Environment=HOME=/home/pi

[Install]
WantedBy=multi-user.target